# 2.1.Linux网关/k7-gateway/systemd/k7-gateway.service
[Unit]
Description=K7 LoRa WiFi MQTT Gateway
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/k7-gateway
Environment=PYTHONPATH=/root/k7-gateway/src
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m k7_gateway run-lora --device /dev/ttyS3 --log /var/log/k7-gateway/lora.jsonl --raw-log /var/log/k7-gateway/lora.raw.log --mqtt-broker broker.emqx.io --transport wifi --mqtt-command-subscribe --quiet
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
